# deploy/monitorkeyword.service
[Unit]
Description=CLS Keyword Monitor (财联社关键词监测)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/root/yangq
ExecStart=/home/root/yangq/.venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
